# systemd/picast.service
[Unit]
Description=Cutegory PiCast Signage Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=picast
Group=picast
SupplementaryGroups=tty video render input
WorkingDirectory=/opt/picast
ExecStart=/opt/picast/picast-client.sh
Restart=always
RestartSec=10
EnvironmentFile=/opt/picast/config.env
TTYPath=/dev/tty7
StandardInput=tty
StandardOutput=journal
StandardError=journal

# Hardening (relaxed for Xorg VT access + log writing)
NoNewPrivileges=false
ProtectSystem=full
ReadWritePaths=/opt/picast /tmp /run /var/log /home/picast
ProtectHome=false

[Install]
WantedBy=multi-user.target
